# deploy/luigi.service
[Unit]
Description=Luigi Health Tracker Bot
After=network.target

[Service]
Type=simple
User=luigi
WorkingDirectory=/opt/luigi_app
Environment=TZ=UTC
EnvironmentFile=/opt/luigi_app/.env
ExecStart=/opt/luigi_app/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luigi

[Install]
WantedBy=multi-user.target
