# deploy/luigi.service
# Prerequisite system packages (install before first start):
#   poppler-utils  — scanned-PDF rasterization (pdf2image)
#   tesseract-ocr  — OCR for image/garbage-layer EOB PDFs (pytesseract)
# Debian/Ubuntu: apt-get install -y poppler-utils tesseract-ocr

[Unit]
Description=Luigi Health Tracker Bot
After=network.target

[Service]
Type=simple
User=luigi
WorkingDirectory=/opt/luigi_app
Environment=TZ=UTC
EnvironmentFile=/opt/luigi_app/.env
ExecStart=/opt/luigi_app/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luigi

[Install]
WantedBy=multi-user.target
